# rc.local
#!/bin/sh -e
/opt/app/start-app.sh
exit 0
